# wsl/ps-admin-shell/profile.d/ps-admin-shell.sh
_PSADMIN_CMD="$HOME/.utils-and-scripts/wsl/ps-admin-shell/main.sh"
alias psadmin="$_PSADMIN_CMD"