# linux/aws-profile-selector/profile.d/aws-profile-selector.sh
#!/bin/bash
# utils-and-scripts/linux/bash-random-generator/profile.d/grand-generator.sh

_ASP_CMD="source $HOME/.utils-and-scripts/linux/aws-profile-selector/main.sh"

alias asp="$_ASP_CMD"
