# linux/bash-random-generator/profile.d/grand-generator.sh
#!/bin/bash
# utils-and-scripts/linux/bash-random-generator/profile.d/grand-generator.sh

_GRAND_CMD="$HOME/.utils-and-scripts/linux/bash-random-generator/main.sh"

alias grand="$_GRAND_CMD"